# services/systemd/edison-core.service
[Unit]
Description=EDISON Core AI Service
After=network.target edison-coral.service
Wants=edison-coral.service

[Service]
Type=simple
User=edison
WorkingDirectory=/opt/edison
Environment="PATH=/opt/edison/.venv/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="CUDA_VISIBLE_DEVICES=0,1,2"
# Runs app.py directly so its __main__ block auto-detects certs for HTTPS
ExecStart=/opt/edison/.venv/bin/python services/edison_core/app.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
