# services/systemd/edison-web.service
[Unit]
Description=EDISON Web UI Service (HTTPS enabled)
After=network.target edison-core.service

[Service]
Type=simple
User=edison
Group=edison
WorkingDirectory=/opt/edison
Environment="PATH=/opt/edison/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/opt/edison/.venv/bin/python services/edison_web/service.py
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
