# services/systemd/edison-web.service
[Unit]
Description=EDISON Web UI Service
After=network.target edison-core.service

[Service]
Type=simple
User=edison
Group=edison
WorkingDirectory=/opt/edison
Environment="PATH=/opt/edison/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
# HTTPS enabled automatically if /opt/edison/certs/cert.pem exists
ExecStart=/opt/edison/.venv/bin/python -m uvicorn services.edison_web.service:app --host 0.0.0.0 --port 8080 --ssl-keyfile /opt/edison/certs/key.pem --ssl-certfile /opt/edison/certs/cert.pem
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
