# services/systemd/edison-comfyui.service
[Unit]
Description=EDISON ComfyUI Image Generation Service
After=network.target
Before=edison-core.service

[Service]
Type=simple
User=edison
WorkingDirectory=/opt/edison/ComfyUI
Environment="PATH=/opt/edison/ComfyUI/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="CUDA_DEVICE_ORDER=PCI_BUS_ID"
ExecStart=/opt/edison/ComfyUI/venv/bin/python main.py --listen 0.0.0.0 --port 8188
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
